# sync_users/crontab
0 0 * * * python /opt/app/sync_users_task.py
